# supervisord.ini
[program:aiserver]
command=/usr/local/aitour/aiserver/aiserver
directory=/usr/local/aitour/aiserver
autostart=true
autorestart=true
startsecs=2
redirect_stderr=true
stdout_logfile=/var/log/aiserver.log


[program:aiweb]
command=/usr/local/aitour/aiweb/aiweb
directory=/usr/local/aitour/aiweb
autostart=true
autorestart=true
startsecs=2
redirect_stderr=true
stdout_logfile=/var/log/aiweb.log



