# service/pythonserver/supervisord.ini
[program:aiserver-py]
command=/home/pangolins/scene/service/pythonserver/start.sh
directory=/home/pangolins/scene/service/pythonserver
autostart=true
autorestart=true
startsecs=2
redirect_stderr=true
stdout_logfile=/var/log/aiserverpy.log


